# todue-restart.service
[Unit]
Description=Restart Todue service
RefuseManualStart=true

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart todue.service
